# sublime-text-3/stage/etc/profile.d/appmenu-qt5.sh
export QT_QPA_PLATFORMTHEME=appmenu-qt5
