# sublime-text-3/parts/desktop/qt5/install/usr/lib/systemd/user/glib-pacrunner.service
[Unit]
Description=GLib proxy auto-configuration service

[Service]
Type=dbus
BusName=org.gtk.GLib.PACRunner
ExecStart=/usr/lib/glib-networking/glib-pacrunner
